# scripts/spotify_jbl_remote_supervisord.conf
[program:spotify_jbl_remote]
command=/usr/bin/spotify-jbl-remote -p
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
user=root